# DOCKER/DOCKER_CLAMAV/COMPOSE/etc/supervisor.d/run.conf
[program:clamspringservice]
command=java -Dspring.profiles.active=%(ENV_SPRING_PROFILE)s -Dlog4j.debug -jar /app.jar
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0